# 02_ft_onion_web/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:tor]
command=/usr/bin/tor -f /etc/tor/torrc
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/sshd.err.log
stdout_logfile=/var/log/sshd.out.log

[program:final_message]
# command=/bin/bash -c 'echo "✅ All processes started successfully!" && echo "🔗 Onion service is available at: http://$(cat /var/lib/tor/hidden_service/hostname)"'
command=/bin/bash -c 'cat /var/lib/tor/hidden_service/hostname'
autostart=true
autorestart=false
startsecs=5
priority=100
